# icpc-wf/ansible/files/create_cgroups.service
[Unit]
Description=Make sure cgroups exist for domjudge judgedaemon

[Service]
Type=oneshot
ExecStart=/opt/domjudge-systest/bin/create_cgroups
RemainAfterExit=true
