# icpc-wf/ansible/files/tune_cpu.service
[Unit]
Description=Tune CPU cores: disable turboboost, downclocking, etc.

[Service]
Type=oneshot
ExecStart=/root/disable_turboboost_ht.sh
RemainAfterExit=true
