# icpc-wf/ansible/files/failover/.config/autostart/alert_listener.desktop
[Desktop Entry]
Name=Alert listener
Type=Application
Exec=bash -c "sleep 10 && ~/bin/alert_listener.sh"
