# icpc-wf/ansible/roles/judgedaemon/templates/domjudge-judgehost.target.j2
[Unit]
Description=DOMjudge judgehost running one or more judgedaemons
Requires={% for core in CPUCORE %}domjudge-judgedaemon@{{ core }}.service {% endfor %}

[Install]
WantedBy=multi-user.target
